# config/systemd/system/omarchy-wayvnc-monitor.service
[Unit]
Description=Omarchy wayvnc Disconnect Lock & Detach Monitor
Documentation=https://github.com/stevepresley/omarchy-advanced
# System service to monitor wayvnc client events and lock+detach on disconnect
# Runs as root to access wayvnc socket owned by root
After=wayvnc.service
Wants=wayvnc.service

[Service]
Type=simple
User=root
# Run the shared monitor script from system location
# This works for any user without duplication
ExecStart=/usr/local/bin/omarchy-wayvnc-monitor
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
